# .config/systemd/user/curl-spam.timer
[Unit]
Description=Run curl-spam every 5 minutes (user timer)

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
