# .config/systemd/user/random-wall.service
[Unit]
Description=Change wallpaper every 20 minutes

[Service]
ExecStart=%h/.local/bin/random_wall.sh
